# linux/stow/systemd/.config/systemd/user/retroarch-saves.service
[Unit]
Description=Backup RetroArch save files to git

[Service]
Type=oneshot
ExecStart=%h/.dotfiles/linux/scripts/backup-retroarch-saves.sh
