# linux/stow/systemd/.config/systemd/user/retroarch-saves.path
[Unit]
Description=Watch RetroArch saves directory for changes
After=local-fs.target

[Path]
PathModified=%h/.dotfiles/shared/symlink/retroarch/.config/retroarch/saves/dolphin-emu/User/GC/USA/Card A
Unit=retroarch-saves.service

[Install]
WantedBy=default.target
